# data/autostart/W10.desktop
[Desktop Entry]
Exec=VBoxManage startvm "Windows 10"
